# docker/jupyterhub_anon/jupyterhub.conf
[supervisord]
nodaemon=true

[program:jupyterhub]
command=/usr/local/bin/jupyterhub -f /srv/jupyterhub/jupyterhub_config.py
autostart = true
autorestart = true
redirect_stderr = true
