# docker/curldrop/curldrop.conf
[supervisord]
nodaemon=true

[program:curldrop]
command=/usr/bin/python /tmp/curldrop-master/runcurldrop.py
autostart = true
autorestart = true
redirect_stderr = true
