# docker/conda_builder/sshd.conf
[program:sshd]
command=bash -c "service sshd start && service sshd stop && exec /usr/sbin/sshd -D -e"
autorestart=true
startretries=3
